# scripts/local-service-booking.service
[Unit]
Description=Local Service Booking Application
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ec2-user/local-service-booking
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
